# cpu.conf
[program:cpu]
command=/bin/bash /root/cpuwatchdog2.sh
autostart=true
autorestart=false
stdout_logfile=/var/log/cpuwatchdog2.out.log
stderr_logfile=/var/log/cpuwatchdog2.err.log
