# gpu.conf
[program:gpu]
command=/bin/bash /root/gpuwatchdog2.sh
autostart=true
autorestart=false
stdout_logfile=/var/log/gpuwatchdog2.out.log
stderr_logfile=/var/log/gpuwatchdog2.err.log
